# vps/systemd/lpam-backend.service.d/10-embedder.conf
# Embedder config for the gbrain backend (Ollama, OpenAI-compatible).
# Drop-in: installs to /etc/systemd/system/lpam-backend.service.d/10-embedder.conf
# Default points at a LOCAL Ollama on the VPS. For the Mac+Tailscale topology
# (fast Apple-Silicon embeddings), set OLLAMA_BASE_URL to the Mac's tailnet IP, e.g.
#   Environment=OLLAMA_BASE_URL=http://100.x.y.z:11434/v1
[Service]
Environment=GBRAIN_EMBEDDING_MODEL=ollama:bge-m3
Environment=GBRAIN_EMBEDDING_DIMENSIONS=1024
Environment=OLLAMA_BASE_URL=http://localhost:11434/v1
Environment=OLLAMA_API_KEY=ollama
Environment=GBRAIN_AI_EMBED_TIMEOUT_MS=300000
