# local/launchagents/com.hermes.ollama-tailnet.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- LaunchAgent: serve Ollama on all interfaces for VPS embeddings over Tailscale.
     TEMPLATE — @@HERMES_BRIDGE@@ substituted by local/bootstrap.sh (default ~/.hermes-bridge). -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.ollama-tailnet</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>@@HERMES_BRIDGE@@/ollama-tailnet.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>@@HERMES_BRIDGE@@/ollama-tailnet.out.log</string>
    <key>StandardErrorPath</key>
    <string>@@HERMES_BRIDGE@@/ollama-tailnet.err.log</string>
</dict>
</plist>
